# files/cron.d/named
SHELL=/bin/bash
2   5 * * * named-update /usr/bin/namedctl cleanup
*/2 * * * * named-update /usr/bin/namedctl update
